# scripts/com.wheelbot.freyrreview.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wheelbot.freyrreview</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3.11</string>
        <string>scripts/freyr_weekly_review.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/openclaw/Documents/btc-wheel-bot</string>
    <!-- Fires DAILY at 5 AM local (AEST) so a sleeping Mac catches up on wake; the
         script itself only writes a NEW review on Sundays (or with --force). -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>5</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/wheelbot_freyrreview.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/wheelbot_freyrreview.log</string>
</dict>
</plist>
